# TURFv6/project-spec/meta-user/recipes-apps/jdld-init/files/jdld-init.service
[Unit]
Description=jdld-init

[Service]
ExecStart=/usr/bin/jdld
StandardOutput=journal

[Install]
WantedBy=multi-user.target
