# TURFv6/project-spec/meta-user/recipes-apps/pueo-squashfs/files/pueo-squashfs.service
[Unit]
Description=pueo-squashfs

[Service]
ExecStart=/usr/bin/pueo-squashfs.sh
StandardOutput=journal
SuccessExitStatus=1
RestartForceExitStatus=0

[Install]
WantedBy=multi-user.target
